# display.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/display/display.sh

[Install]
WantedBy=default.target

